# sm_config/bar.service
[Unit]
Description=bar

[Service]
ExecStart=LOG=".log" /usr/bin/tail -f /var/log/system.log
Restart=always

[Install]
WantedBy=default.target

