# sm_config/client5.service
[Service]
ExecStart=/usr/bin/java -Dtimes=5 -jar /Users/renren/app/client5/target/client5-1.0.jar
WorkingDirectory=/Users/renren/sm
