# sm_config/hi.service
[Service]
ExecStart=/usr/local/bin/node /Users/renren/app/hi.js
