# sm_config/sm-monitor.service
[Service]
ExecStart=/usr/bin/sm monitor
Restart=always
User=renren

[Install]
WantedBy=default.target
