# sm_config/tb.service
[Service]
ExecStart=/usr/bin/java -Dserver.port=9414 -Xms256M -Xmx256M -jar /Users/renren/sm/turbine/target/turbine-2.0.jar
